# service-definitions/systemd/lxdepot.service
[Unit]
Description=LXDepot
After=network.target

[Service]
WorkingDirectory=/opt/lxdepot/
ExecStart=/opt/lxdepot/lxdepot
Restart=always
